# systemd/leovisa-frontend.service
[Unit]
Description=LeoVisa Digital Human Frontend
After=network.target leovisa-backend.service

[Service]
Type=simple
WorkingDirectory=/home/ryanrenjr/AI-Workspace/app/frontend
ExecStartPre=/bin/mkdir -p /home/ryanrenjr/AI-Workspace/logs
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=3
StandardOutput=append:/home/ryanrenjr/AI-Workspace/logs/frontend.log
StandardError=append:/home/ryanrenjr/AI-Workspace/logs/frontend.log

[Install]
WantedBy=default.target
