# systemd/leovisa-backend.service
[Unit]
Description=LeoVisa Digital Human Backend
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
WorkingDirectory=/home/ryanrenjr/AI-Workspace/app/backend
ExecStartPre=/bin/mkdir -p /home/ryanrenjr/AI-Workspace/logs
ExecStart=/home/ryanrenjr/miniconda3/bin/python -m uvicorn main:app --host 127.0.0.1 --port 8008
Restart=always
RestartSec=3
StandardOutput=append:/home/ryanrenjr/AI-Workspace/logs/backend.log
StandardError=append:/home/ryanrenjr/AI-Workspace/logs/backend.log

[Install]
WantedBy=default.target
